# etc/systemd/system/batchjob_stats_server.service
[Unit]
Description=Batchjob Stats API server
After=network-online.target

[Service]
Type=simple
User=batchjob-stats
Group=batchjob-stats
ExecStart=/usr/local/bin/batchjob_stats_server \
    --slurm.sacct.path="/usr/local/bin/sacct" \
    --slurm.sacct.run.with.sudo \
    --path.data="/var/lib/batchjob_stats" \
    --db.update.interval=1800 \
    --data.retention.period=365 \
    --log.level="debug"

SyslogIdentifier=batchjob_stats_server
Restart=always
RestartSec=1
StartLimitInterval=0

NoNewPrivileges=yes

ProtectSystem=strict
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
